# deploy/systemd/tcp-over-kafka-client.service
[Unit]
Description=TCP over Kafka SOCKS5 client proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/libexec/tcp-over-kafka/run-client.sh
Restart=always
RestartSec=1
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
