# deploy/systemd/tcp-over-kafka-server.service
[Unit]
Description=TCP over Kafka relay server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/libexec/tcp-over-kafka/run-server.sh
Restart=always
RestartSec=1
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
